# etc/supervisor.d/websockets.ini
[program:ws]
command = php /opt/ows/artisan websockets:serve
autostart = true
autorestart = true
numprocs = 1
