# etc/supervisor.d/websockets.ini
[program:ws]
command = php /opt/ows/artisan websockets:serve
autostart = true
autorestart = true
numprocs = 1
user=root
redirect_stderr = true
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
stderr_logfile_maxbytes = 0